# throttle.service
[Unit]
Description=Temperature-based CPU throttling

[Service]
ExecStart=PREFIX/sbin/throttle CONFIG_FILE PIPE

[Install]
WantedBy=multi-user.target
